# 02_recursos/feveem.socket
[Unit]
Description=Socket de la API de Censo Estudiantil

[Socket]
ListenStream=/run/feveem.sock

[Install]
WantedBy=sockets.target